# deploy/mission-control-scraper.service
[Unit]
Description=Mission Control — scraper task (runs once per timer trigger)
After=network.target

[Service]
Type=exec
WorkingDirectory=/opt/mission-control
ExecStart=/usr/local/bin/bun run src/workers/scraper-worker.ts
Environment=NODE_ENV=production

[Install]
WantedBy=timers.target
