# deploy/mission-control-energy-price-scraper.timer
[Unit]
Description=Run energy-price scraper daily at 8:00 AM

[Timer]
OnCalendar=*-*-* 08:00:00
Persistent=true
RandomizedDelaySec=120

[Install]
WantedBy=timers.target
